# docker/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/nginx.log
stderr_logfile=/var/log/supervisor/nginx_err.log

[program:api]
command=/bin/sh -c 'export PYTHONPATH=/app PYTHONUNBUFFERED=1 && exec uvicorn backend.main:app --host 0.0.0.0 --port 8000'
directory=/app
autostart=true
autorestart=true
stdout_logfile=/app/logs/api.log
stderr_logfile=/app/logs/api_err.log

[program:worker]
command=/bin/sh -c 'export PYTHONPATH=/app PYTHONUNBUFFERED=1 && exec python -m worker.main'
directory=/app
autostart=true
autorestart=true
stdout_logfile=/app/logs/worker.log
stderr_logfile=/app/logs/worker_err.log
stopwaitsecs=10
